# files/system/etc/systemd/system/uupd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/uupd --disable-module-distrobox --disable-module-brew --log-level debug --json --hw-check
